# ops/astro-blog-pull-deploy.service
[Unit]
Description=Pull source, build, and deploy astro-blog
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
EnvironmentFile=-/etc/astro-blog-pull-deploy.env
ExecStart=/usr/local/bin/astro-blog-pull-deploy
TimeoutStartSec=30min
