# ops/astro-blog-pull-deploy.timer
[Unit]
Description=Check for astro-blog deploy artifacts every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=astro-blog-pull-deploy.service

[Install]
WantedBy=timers.target
